# etc/systemd/system/gwmp_mux.service
[Unit]
Description=GWMP UDP Multiplexer (LoRa Packet Forwarder → TTN + local consumer)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ttnuser
Group=ttnuser

WorkingDirectory=/home/ttnuser/ttn-message-interceptor
ExecStart=/home/ttnuser/ttn-message-interceptor/gwmp_mux.py

Environment=MUX_LISTEN_PORT=1700
Environment=MUX_UPSTREAM_HOST=eu1.cloud.thethings.network
Environment=MUX_UPSTREAM_PORT=1700
Environment=MUX_CONSUMER_HOST=127.0.0.1
Environment=MUX_CONSUMER_PORT=1701
Environment=MUX_LOG_LEVEL=INFO

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=gwmp-mux

[Install]
WantedBy=multi-user.target
